# .config/systemd/user/claude-web-ui.service
[Unit]
Description=Claude Web UI
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/home/zach/Source/Projects/claude-web-ui
ExecStartPre=/bin/bash -c "./build.sh"
ExecStart=/usr/sbin/stack exec claude-web-ui 
Restart=always

[Install]
WantedBy=default.target
