# .config/systemd/user/dbox-start-util.service
[Service]
Type=oneshot
RemainAfterExit=true
StandardOutput=journal
ExecStart=/usr/bin/distrobox-enter util -- echo started

[Install]
WantedBy=default.target

