# .config/systemd/user/mount-drives.service
[Service]
Type=oneshot
RemainAfterExit=true
StandardOutput=journal
ExecStart=/usr/bin/sudo /usr/bin/mount -a

[Install]
WantedBy=default.target
