# .config/systemd/user/nas-sleep-guard.service
[Unit]
Description=Unmount NAS sshfs before suspend (avoid D-state suspend-freeze hang)
Documentation=man:systemd-inhibit(1)
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.dotfiles/bin/scripts/nas-sleep-guard
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
