# .config/systemd/user/var-home-zach-data-mnt-nas-apool.mount
[Unit]
Description=SSHFS Mount for NAS (apool)
After=network-online.target tailscaled.service
Wants=network-online.target tailscaled.service
# Retry configuration
StartLimitIntervalSec=300
StartLimitBurst=5

# Make sure user is part of fuse group 
# or `user_allow_other` is uncommented in /etc/fuse.conf
[Mount]
What=zach@nas-main:/var/mnt/apool
# systemd-escape -p --suffix=mount "/home/zach/.mnt/nas"
Where=%h/data/mnt/nas/apool
Type=fuse.sshfs
# --- Options changed 2026-06-08 (suspend-hang fix) ----------------------------
# `reconnect` makes sshfs retry a dead link forever, pinning any task doing I/O
# on this mount in uninterruptible (D) state. The suspend freezer cannot freeze
# a D-state task, so suspend aborts ("N tasks refusing to freeze") and the
# laptop hangs needing a hard power-off. Started with the 7.0.11 kernel
# (Immutablue 44.20260606.0). Fix: drop `reconnect` (a dead link now returns EIO
# instead of blocking forever), tighten the SSH keepalive so a dead link is
# detected (~10s) inside the 20s freeze window. NOTE: x-systemd.automount and
# x-systemd.idle-timeout are INERT in a native .mount unit (they are only
# honored for /etc/fstab entries; no companion .automount unit exists), so the
# mount does not auto-unmount when idle -- the active wins here are dropping
# `reconnect` and the tighter keepalive. Resume remounts are handled by the
# nas-sleep-guard.service. To revert, re-enable the original line below.
# Original:
# Options=_netdev,reconnect,ServerAliveInterval=15,ServerAliveCountMax=3,x-systemd.automount,IdentityFile=%h/.ssh/id_rsa,allow_other,default_permissions,ConnectTimeout=10,compression=no
Options=_netdev,ServerAliveInterval=5,ServerAliveCountMax=2,x-systemd.automount,x-systemd.idle-timeout=120,IdentityFile=%h/.ssh/id_rsa,allow_other,default_permissions,ConnectTimeout=10,compression=no
TimeoutSec=90

[Install]
WantedBy=default.target

