# .config/autostart/ch.protonmail.protonmail-bridge.desktop
[Desktop Entry]
Type=Application
Name=Proton Mail Bridge
Exec=flatpak run ch.protonmail.protonmail-bridge
X-Flatpak=ch.protonmail.protonmail-bridge
