# .config/containers/systemd/ollama-libreclaw-00.container
[Unit]
Description=Ollama container
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=0

[Container]
Image=docker.io/ollama/ollama:latest
PublishPort=127.0.0.1:11434:11434
PublishPort=100.123.41.15:11434:11434
Volume=/var/home/zach/.ollama:/root/.ollama:z

# Number of parallel requests to handle simultaneously (default: 4 or 1 based on memory)
Environment=OLLAMA_NUM_PARALLEL=4

# Maximum number of models to keep loaded simultaneously (default: 1)
Environment=OLLAMA_MAX_LOADED_MODELS=2

# Enable Flash Attention optimization for better memory efficiency
Environment=OLLAMA_FLASH_ATTENTION=1

[Service]
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
